# web/crontab
*/5 * * * * root bash /app/run-background-worker.sh >> /app/cron.log 2>&1
